# crontab
#Enable automatic certificate renewal
0 5 * * * /usr/bin/certbot renew --quiet